# supervisor.conf
[program:kibana4]
user=root
command=/kibana/run.sh
redirect_stderr=true
autorestart=true
